# Library/LaunchAgents/com.gyndok.clipboard-appender.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gyndok.clipboard-appender</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
    <string>/Users/gyndok/clipboard-appender/clipboard_appender.py</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/gyndok/clipboard-appender</string>

  <key>StandardOutPath</key>
  <string>/Users/gyndok/.config/clipboard-appender/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/gyndok/.config/clipboard-appender/launchd.err.log</string>
</dict>
</plist>
